# meta-aos-rcar-common/meta-aos-rcar-common-domx/recipes-aos/aos-updatemanager/files/aos-reboot.service
[Unit]
Description=AOS Reboot service
ConditionPathExists=/var/aos/workdirs/um/update_rootfs/updated

[Service]
ExecStart=systemctl reboot --force

[Install]
WantedBy=multi-user.target
