# meta-aos-rcar-common/meta-aos-rcar-common-domd/recipes-aos/aos-updatemanager/files/aos-reboot.service
[Unit]
Description=AOS Reboot service
ConditionPathExists=/var/aos/workdirs/um/update_rootfs/updated

[Service]
ExecStart=/usr/bin/xenstore-write control/user-reboot 2

[Install]
WantedBy=multi-user.target
