# meta-aos-rcar-gen3/meta-aos-rcar-gen3-domd/recipes-aos/aos-updatemanager/files/aos-reboot.service
[Unit]
Description=AOS Reboot service
ConditionPathExists=/var/aos/downloads/update_rootfs/updated

[Service]
ExecStart=systemctl reboot --force

[Install]
WantedBy=multi-user.target
