# application_cpu/meta-rcar-whitebox/meta-rcar-common/recipes-sysmonitor/node-exporter/files/node_exporter.service
[Unit]
Description = Node Exporter
After=local-fs.target
ConditionPathExists=/usr/bin

[Service]
ExecStart=/usr/bin/node_exporter/node_exporter
WorkingDirectory=/usr/bin/node_exporter
Restart=no
Type=simple

[Install]
WantedBy=multi-user.target

