# application_cpu/meta-rcar-whitebox/meta-rcar-common/recipes-connectivity/covesa-ccs/files/agt_server.service
[Unit]
Description = WAII : agt_server
After=local-fs.target network.target
PartOf=waii.target
ConditionPathExists=/usr/bin

[Service]
ExecStart=/usr/bin/waii/agt_server/agt_server
WorkingDirectory=/usr/bin/waii/agt_server
Type=simple
Restart=on-failure
RestartSec=10

[Install]
#WantedBy=multi-user.target

