# application_cpu/meta-rcar-whitebox/meta-rcar-common/recipes-connectivity/covesa-ccs/files/waii.target
[Unit]
Description = WAII control service file

[Unit]
Wants=agt_server.service
Wants=vissv2server.service

[Install]
WantedBy=multi-user.target

