# meta-aos-rcar-gen4/meta-aos-rcar-gen4-domd/recipes-bsp/nvme-cli/files/nvme.service
[Unit]
Description=nvme service
ConditionFileNotEmpty=!/etc/nvme/hostnqn

[Service]
Type=oneshot
ExecStart=sh -c "nvme gen-hostnqn > /etc/nvme/hostnqn"
ExecStart=sh -c "uuidgen > /etc/nvme/hostid"

[Install]
WantedBy=multi-user.target
