# meta-aos-rcar-common/meta-aos-rcar-common-dom0/recipes-guests/guestreboot/files/guestreboot@.service
[Unit]
Description=Reboot system by guest domain request

[Service]
ExecStart=/usr/bin/guestreboot.sh %i
